# NEW/puppet-cobbler/files/cobbler-replicate.cron
#!/bin/bash

# Dynamically generated by Puppet
# Local modifications will be over-written

# Replicate data to slave, just in case
# Is also handled via Cobbler tiggers in /var/lib/cobbler/triggers
/usr/local/bin/cobbler-replicate &>> /dev/null
