# NEW/puppet-cobbler/files/cobbler-reposync.cron
#!/bin/bash

# Dynamically generated by Puppet
# Local modifications will be over-written

# Clean out tasks log directory
# (not a great candidate for logrotate)
/bin/find /var/log/cobbler/tasks -maxdepth 1 -mmin +10080 -name "*.log" -exec rm -f {} \;

#run reposync to update repos
/usr/bin/cobbler reposync --tries=3 --no-fail &>> /dev/null
